# actn/fleet/instances/cors-delete-discovery@devos-1.service
[Unit]
Description=Actn.io %p %i service discovery

# Requirements
Requires=etcd.service
Requires=docker.service
Requires=cors-delete@%i.service

# Dependency ordering and binding
After=etcd.service
After=docker.service
After=cors-delete@%i.service
BindsTo=cors-delete@%i.service
#Before=%p-notify@%i.service


[Service]
EnvironmentFile=/etc/environment
ExecStart=/bin/sh -c "while true; do etcdctl set /apps/devos/apis/cors/delete/%i $(docker port $(docker ps | grep cors-delete-%i | awk '{print $1}') 9000) --ttl 30; sleep 20; done"
ExecStop=/usr/bin/etcdctl rm /apps/devos/apis/cors/delete/%i

#[Install]
#WantedBy=app@devos.target

[X-Fleet]
X-ConditionMachineOf=cors-delete@%i.service
